# resources/systemd/porter@.service
[Unit]
Description=Porter workspace: %i
Documentation=https://github.com/dougefresher/pi-porter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/%i
EnvironmentFile=%h/%i/.env
# Prod unit: packaged install at a fixed path. PATH only affects subprocesses, not ExecStart.
ExecStart=/usr/bin/porter --serve
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
RuntimeDirectory=porter
RuntimeDirectoryMode=0700
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
RestartSteps=4
RestartMaxDelaySec=160s
KillSignal=SIGTERM
KillMode=mixed

[Install]
WantedBy=multi-user.target
